# com.anyrouter.checkin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anyrouter.checkin</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>checkin.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/robert/Projects/anyrouter</string>

    <!-- 每天 08:00 执行，睡眠/关机唤醒后自动补跑 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/robert/Projects/anyrouter/checkin.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/robert/Projects/anyrouter/checkin.log</string>
</dict>
</plist>
